# supervisord.conf
# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:listener]
command=node listener.js
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0```

**4. Update `listener.js` to Trigger Render's Deploy Hook**

This is the final magic. Instead of running `git pull`, the listener will make a simple web request to a special URL provided by Render, telling it to redeploy the service.

```javascript
// listener.js (Updated for Render)
import express from 'express';
import crypto from 'crypto';

const app = express();
const PORT = 9001; // Internal port, not exposed to the public
const WEBHOOK_SECRET = process.env.WEBHOOK_SECRET;
const RENDER_DEPLOY_HOOK_URL = process.env.RENDER_DEPLOY_HOOK_URL;

if (!WEBHOOK_SECRET || !RENDER_DEPLOY_HOOK_URL) {
  console.error("FATAL: WEBHOOK_SECRET and RENDER_DEPLOY_HOOK_URL must be set.");
  process.exit(1);
}

app.use(express.json());

app.post('/webhook', (req, res) => {
  console.log('Webhook received...');
  // Security check remains the same...
  const signature = req.headers['x-hub-signature-256'];
  if (!signature) return res.status(401).send('Signature required.');
  const hmac = crypto.createHmac('sha256', WEBHOOK_SECRET);
  const digest = 'sha256=' + hmac.update(JSON.stringify(req.body)).digest('hex');
  if (!crypto.timingSafeEqual(Buffer.from(digest), Buffer.from(signature))) {
    return res.status(401).send('Invalid signature.');
  }

  console.log('Signature verified. Triggering Render deploy hook...');
  
  // THE NEW LOGIC: Instead of exec, we call the deploy hook.
  fetch(RENDER_DEPLOY_HOOK_URL)
    .then(hookRes => {
      if (hookRes.ok) {
        console.log('✅ Successfully triggered Render deployment.');
        res.status(200).send('Deployment triggered successfully!');
      } else {
        console.error(`Failed to trigger Render hook: ${hookRes.statusText}`);
        res.status(500).send('Failed to trigger deployment.');
      }
    })
    .catch(err => {
      console.error('Error calling Render deploy hook:', err);
      res.status(500).send('Error triggering deployment.');
    });
});

app.listen(PORT, '127.0.0.1', () => {
  console.log(`Internal webhook listener started on port ${PORT}`);
});
